# supervisord.conf

[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:run]
command=/bin/sh -c "/usr/local/bin/run"



